# platform-runner.service
[Unit]
Description=GitHub Actions runners (grid)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Replace <USER> with the host user that owns the grid checkout,
# and <PATH> with the absolute path to that checkout.
User=<USER>
WorkingDirectory=<PATH>
ExecStart=/bin/bash -c "docker compose up --build --scale runner=4 --detach"
ExecStop=/bin/bash -c "docker compose stop"

[Install]
WantedBy=multi-user.target
